# salt/roots/salt/ipcluster/etc/supervisor/conf.d/ipengine.conf
{% set username = pillar.get('ipcluster.username', 'ipuser') %}
{% set home = pillar.get('ipcluster.userhome', '/home/ipuser') %}
{% set host_cpus = salt['ps.num_cpus']() %}
{% set controller_cpus = pillar.get('ipcluster.controller.cpus', 0) %}
{% set engine_cpus = host_cpus - controller_cpus %}

[program:ipengine]
command=ipcluster engines -n {{ engine_cpus }}
user={{ username }}
directory={{ pillar.get('ipcluster.engine.directory', home) }}
environment=HOME='{{ home }}',USER={{ username }}
