# salt/roots/salt/ipcluster/etc/supervisor/conf.d/ipengine.conf
# To be deployed on worker-only node that do not have an instance of the
# IPython controller
{% set host_cpus = salt['ps.num_cpus']() %}

# TODO: add the SSH config to the host controller

[program:ipengine]
command={{ ipcluster }} engines -n {{ host_cpus }}
user={{ user }}
directory={{ directory }}
environment=HOME='{{ home }}',USER={{ user }}
stopwaitsecs=10
# ipcluster does not propagate the TERM signal to children by default
stopasgroup=true