# salt/roots/salt/ipcluster/etc/supervisor/conf.d/ipengine.conf
[program:ipengine]
command=ipcluster engines
process_name=ipcluster-engines
user=root

# TODO: use templating for using a dynamic, non-priviledged user and set
# command arguments