# salt/roots/salt/ipcluster/etc/supervisor/conf.d/ipnotebook.conf
[program:ipnotebook]
command=ipnotebook
process_name=ipnotebook
user=root

# TODO: use templating for using a dynamic, non-priviledged user and set
# command arguments