# salt/roots/salt/ipcluster/etc/supervisor/conf.d/ipnotebook.conf
{% set username = pillar.get('ipcluster.username', 'ipuser') %}
{% set home = pillar.get('ipcluster.userhome', '/home/ipuser') %}

[program:ipnotebook]
command=ipython notebook
    --profile={{ pillar.get('ipcluster.profile', 'default') }}
    --port={{ pillar.get('ipcluster.notebook.port', '8888') }}
    --ip=*
    --no-browser
user={{ username }}
autostart={{ pillar.get('ipcluster.notebook.enabled', True) }}
directory={{ pillar.get('ipcluster.notebook.directory', home) }}
environment=HOME='{{ home }}',USER={{ username }}
